# amanda-srv.crontab
0 16 * * 1-5	amanda	/usr/sbin/amcheck -m
45 0 * * 2-6	amanda	/usr/sbin/amdump
